# .github/workflows/farmr.service
[Unit]
Description=farmr

[Service]
Type=simple
WorkingDirectory=/etc/farmr/
ExecStart=/usr/bin/farmr

[Install]
WantedBy=default.target